# crontab
*/10 * * * * /home/lsund/Documents/dotfiles/shell/memlog >/dev/null 2>&1
*/15 * * * * /home/lsund/Documents/dotfiles/shell/pacsync >/dev/null 2>&1
